# edge/sync/systemd/eunomia-sweep-broken.timer
[Unit]
Description=Eunomia SYNC1 broken-episode sweep (hourly)

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
AccuracySec=5min
Persistent=true
Unit=eunomia-sweep-broken.service

[Install]
WantedBy=timers.target
